# files/docker-prune.timer
[Unit]
Description=Run Docker prune weekly

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
